# debian/cm4all-qrelay.service
[Unit]
Description=CM4all QMQP relay daemon
Wants=network-online.target
After=network-online.target

[Service]
Type=notify
User=cm4all-qrelay
ExecStart=/usr/sbin/cm4all-qrelay --config /etc/cm4all/qrelay/config.lua

CPUSchedulingPolicy=batch

# enable crash dumps
LimitCORE=infinity

# Paranoid security settings
NoNewPrivileges=yes
ProtectSystem=yes
ProtectHome=yes
ProtectKernelTunables=yes
ProtectControlGroups=yes
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=yes
