# debian/cm4all-qrelay.service
[Unit]
Description=CM4all QMQP relay daemon
After=network.target

[Service]
Environment=CONFIG=/etc/cm4all/qrelay.conf
EnvironmentFile=/etc/default/cm4all-qrelay
ExecStart=/usr/sbin/cm4all-qrelay --no-daemon --config ${CONFIG} --user cm4all-qrelay

[Install]
WantedBy=multi-user.target
