# debian/cm4all-qrelay.socket
[Socket]
ListenStream=/run/cm4all/qrelay/socket
PassCredentials=true

[Install]
WantedBy=sockets.target
